# scripts/lms.service
[Unit]
Description=LMS

[Service]
Type=forking
ExecStart=/usr/bin/lms

[Install]
WantedBy=multi-user.target
